# systemd/grott-scheduler.service
[Unit]
Description=Grott Scheduler Service
Documentation=https://github.com/johanmeijer/grott
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/grott-scheduler/backend
ExecStart=/usr/bin/python3 /opt/grott-scheduler/backend/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grott-scheduler

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
